# etc/NetworkManager/dispatcher.d/99-tailscale-ondemand
#!/bin/bash
# Trigger the script when a network connects or disconnects
if [ "$2" = "up" ] || [ "$2" = "down" ]; then
    /usr/local/bin/ts-ondemand
fi
